# replay-ws.service
[Service]
ExecStart=/usr/bin/python3 /home/pi/replay/websocket/server.py
